# TruckSimulatorBot.service
[Unit]
Description=Truck Simulator Discord bot
After=network-online.target
After=mariadb.service

[Service]
User=USER
WorkingDirectory=WORKINGDIRECTORY
ExecStart=WORKINGDIRECTORY/venv/bin/python3 -m gunicorn bot:app -b 127.0.0.1:9000 -w 2

[Install]
WantedBy=network-online.target

