# TruckSimulatorBot.service
[Unit]
Description=Truck Simulator Discord bot

[Service]
User=USER
WorkingDirectory=WORKINGDIRECTORY
ExecStart=WORKINGDIRECTORY/venv/bin/python3 bot.py --enable-log-file --post-server-count

[Install]
WantedBy=multi-user.target
